# templates/systemd/subscription-aggregator.service
[Unit]
Description=Reality residential subscription aggregator service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/reality-resi-stack/subscription-aggregator.env
ExecStart=/usr/bin/python3 /usr/local/lib/reality-resi-stack/aggregator_server.py
Restart=always
RestartSec=3
WorkingDirectory=/usr/local/lib/reality-resi-stack

[Install]
WantedBy=multi-user.target
